# linux/roles/wod.flannel/templates/flannel-docker-opts.service.d/10-image.conf
[Service]                                     
Environment="FLANNEL_IMAGE=/etc/kubernetes/downloads/{{ RKT_ACI_FLANNEL }}.aci"    
Environment="RKT_GLOBAL_ARGS=--insecure-options=image"
